# build_files/kde-builder-session.service
[Unit]
Description=Install KDE builder login session unless user overrides it
DefaultDependencies=no
After=local-fs.target
Before=display-manager.service

[Service]
Type=oneshot
ExecStart=/usr/bin/kde-builder-session-guard.sh

[Install]
WantedBy=multi-user.target
